# orangepi-build/external/packages/bsp/common/lib/systemd/system/setup-overlay.service
[Unit]
Description=Setup overlayfs automatically once
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup_overlay.sh

[Install]
WantedBy=multi-user.target

