# orangepi-build/external/packages/bsp/overlays_cix/etc/profile.d/cix_env.sh
export KWIN_COMPOSE=O2ES
export QT_OPENGL=es2
export GDK_GL=gles
export GST_GL_API=gles2
